# packages/lxde-desktop/lxdm/files/systemd/lxdm.service
[Unit]
Description=Display Manager
Conflicts=kdm.service gdm.service slim.service xdm.service

[Service]
ExecStart=/usr/sbin/lxdm-binary

[Install]
Alias=displaymanager.service
WantedBy=graphical.target
